# system/systemd/nvidia_fan_speed.timer
[Unit]
Description=Set NVIDIA fan speed

[Timer]
OnBootSec=5s
OnUnitActiveSec=4s
AccuracySec=1s

[Install]
WantedBy=timers.target
